# roles/open_webui/templates/open-webui-brain-sync.service.j2
{{ ansible_managed | comment }}
# One brain-sync cycle. Fired by open-webui-brain-sync.timer; runs as ROOT so
# it can rewrite /etc/open-webui.env and restart open-webui.
[Unit]
Description=Open WebUI brain runtime sync (repoint the default model from OpenBao, no rebuild)
After=network-online.target open-webui.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-{{ open_webui_brain_sync_env_path }}
ExecStart={{ open_webui_brain_sync_script_path }}
