# roles/open_webui/templates/open-webui-brain-sync.timer.j2
{{ ansible_managed | comment }}
# Drive the Open WebUI brain runtime sync every {{ open_webui_brain_sync_interval }}.
[Unit]
Description=Periodic Open WebUI brain runtime sync (OpenBao)

[Timer]
OnBootSec={{ open_webui_brain_sync_interval }}
OnUnitActiveSec={{ open_webui_brain_sync_interval }}
AccuracySec=5s

[Install]
WantedBy=timers.target
