# roles/ai_runner/templates/ai-job-worker.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=AI job worker (polls the Vikunja AI Jobs queue, runs coding-agent jobs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ai_runner_user }}
Group={{ ai_runner_group }}
WorkingDirectory={{ ai_runner_home }}
EnvironmentFile={{ ai_runner_worker_env_file }}
ExecStart=/usr/bin/python3 {{ ai_runner_worker_script }}
# The worker is a 24/7 poll loop — any exit is an outage, so restart always.
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
