# roles/hermes_agent/templates/hermes-brain-sync.service.j2
{{ ansible_managed | comment }}
# One brain-sync cycle. Fired by hermes-brain-sync.timer; runs as ROOT
# (separate from the hermes user's gateway process — see defaults/main.yml)
# so it can rewrite config.yaml and restart hermes-gateway.
[Unit]
Description=Hermes brain runtime sync (repoint the fabric brain from OpenBao, no rebuild)
After=network-online.target hermes-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-{{ hermes_agent_brain_sync_env_path }}
ExecStart={{ hermes_agent_brain_sync_script_path }}
