# roles/hermes_agent/templates/hermes-brain-watchdog-alert.service.j2
{{ ansible_managed | comment }}
# OnFailure= target for hermes-brain-watchdog.service: fires only when a probe
# cycle CRASHES (unhandled error), never for a detected brain-down edge (those
# exit 0 and are alerted separately). Rate-limited via StartLimit* below so a
# persistently crashing watchdog pages once per window instead of once per
# {{ hermes_agent_brain_watchdog_interval }} cycle -- systemd enforces this
# start-rate limit for ALL start attempts of this unit, including ones
# triggered by another unit's OnFailure=, not just Restart= loops.
#
# Path + thresholds are deliberately literal, not role variables: they are
# operational constants tied 1:1 to this alert unit, and nothing overrides
# them per-host today. Keep the ExecStart= path in sync with the `dest:` in
# tasks/main.yml's "Deploy the Hermes brain-health watchdog self-check alert
# script" task if either one ever moves.
[Unit]
Description=Alert that the Hermes brain-health watchdog itself failed
StartLimitIntervalSec=1h
StartLimitBurst=1

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hermes-brain-watchdog-alert.sh
