# roles/hermes_agent/templates/hermes-brain-sync.timer.j2
{{ ansible_managed | comment }}
# Drive the Hermes brain runtime sync every {{ hermes_agent_brain_sync_interval }}.
[Unit]
Description=Periodic Hermes brain runtime sync (OpenBao)

[Timer]
OnBootSec={{ hermes_agent_brain_sync_interval }}
OnUnitActiveSec={{ hermes_agent_brain_sync_interval }}
AccuracySec=5s

[Install]
WantedBy=timers.target
