# roles/hermes_agent/templates/hermes-dashboard.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Hermes Dashboard (authenticated interactive UI)
After=network-online.target hermes-gateway.service
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_agent_user }}
Group={{ hermes_agent_user }}
WorkingDirectory={{ hermes_agent_hermes_home }}
Environment=HERMES_HOME={{ hermes_agent_hermes_home }}
EnvironmentFile=-{{ hermes_agent_hermes_home }}/.env
ExecStart={{ hermes_agent_bin }} dashboard --host {{ hermes_agent_dashboard_host }} --port {{ hermes_agent_dashboard_port }} --no-open
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
